# firewall-nftables/etc/init.d/firewall-nft
#!/sbin/runscript
#######     Firewall escrito por Raphael Bastos  #############
###########  raphaelbastos@hackstore.com,br  #################
#      Este script foi escrito para a seguinte distro:       #
#      * Funtoo Linux *                                      #
##############################################################

depend() {
	need net
	need nftables
	}

start(){
##################### PARTE 1 #####################
ebegin "Iniciando o Firewall NFTables"

# Cria link simbólico RCFIREWALL caso não exista
if [ ! -L /usr/sbin/rcfirewall-nft ];then
	ln -s /etc/init.d/firewall-nft /usr/sbin/rcfirewall-nft 2> /dev/null
fi

#######################################################################
### Incompatibilities ###
# You cannot use iptables and nft to perform NAT at the same time. So make sure that the iptable_nat module is unloaded
rmmod iptable_nat

#######################################################################

echo 1 > /proc/sys/net/ipv4/ip_forward
echo 1 > /proc/sys/net/ipv4/ip_dynaddr
echo 1 > /proc/sys/net/ipv4/conf/all/rp_filter
for f in /proc/sys/net/ipv4/conf/*/rp_filter ; do echo 1 > $f ; done

#######################################################################

iptables -t nat -F

#######################################################################

# ipv4
nft -f /etc/nftables/ipv4-filter

# ipv4 nat
nft -f /etc/nftables/ipv4-nat

# ipv6
nft -f /etc/nftables/ipv6-filter

# regras de firewall NTFtables
nft -f /etc/nftables/firewall.rules

#######################################################################

}

stop(){

#######################################################################

#iptables -t nat -F
NFT=nft
FAMILIES="ip ip6 arp bridge"

for FAMILY in $FAMILIES; do
  TABLES=$($NFT list tables $FAMILY | grep "^table\s" | cut -d' ' -f2)
  
  for TABLE in $TABLES; do
    CHAINS=$($NFT list table $FAMILY $TABLE | grep "^\schain\s" | cut -d' ' -f2)

    for CHAIN in $CHAINS; do
      echo "Flushing chain: $FAMILY->$TABLE->$CHAIN"
      $NFT flush chain $FAMILY $TABLE $CHAIN
      $NFT delete chain $FAMILY $TABLE $CHAIN
    done

    echo "Flushing table: $FAMILY->$TABLE"
    $NFT flush table $FAMILY $TABLE
    $NFT delete table $FAMILY $TABLE
  done
done
}

status(){
nft list ruleset
}

# Fim do script
