# firewall-nftables/etc/init.d/firewall-nft
#!/sbin/runscript
#######     Firewall escrito por Raphael Bastos  #############
###########  raphaelbastos@hackstore.com,br  #################
#      Este script foi escrito para a seguinte distro:       #
#      * Funtoo Linux *                                      #
##############################################################

depend() {
	need net
	}

start(){
##################### PARTE 1 #####################
ebegin "Iniciando o Firewall NFTables"

# Cria link simbólico RCFIREWALL caso não exista
if [ ! -L /usr/sbin/rcfirewall-nft ];then
	ln -s /etc/init.d/firewall-nft /usr/sbin/rcfirewall-nft 2> /dev/null
fi

#######################################################################

# ipv4
nft -f /etc/nftables/ipv4-filter

# ipv4 nat
nft -f /etc/nftables/ipv4-nat

# ipv6
nft -f /etc/nftables/ipv6-filter

# regras de firewall NTFtables
nft -f /etc/nftables/firewall.rules

}

stop(){
# ipv4
nft delete rule filter input
nft delete rule filter output
nft delete rule filter forward

# ipv6
nft delete rule ip6 filter input
nft delete rule ip6 filter output
nft delete rule ip6 filter forward

# ipv4 nat
nft delete rule ip nat input
nft delete rule ip nat prerouting
nft delete rule ip nat postrouting
}

status(){
nft list ruleset
}

# Fim do script
